# infra/minimal/supervisord.conf
; supervisord — runs every Aegis process inside the single aegis-core container.
;
; Each Python service binds 127.0.0.1:<port>. Only the gateway is reachable
; from outside the container (it binds 0.0.0.0:8000). All inter-service
; calls are loopback. A CISO inspecting the container sees: one external
; port, one set of in-process services, one security perimeter.

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ── OPA — runs first, listens on 127.0.0.1:8181 ─────────────────────────
; Started before the policy service so the bundle loads have an OPA to
; talk to. `--addr 127.0.0.1:8181` ensures it's not reachable from
; outside the container.
[program:opa]
command=/usr/local/bin/opa run --server --addr 127.0.0.1:8181 /app/services/policy/policies
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── (PgBouncer omitted in minimal mode) ────────────────────────────────
; The full prod-ha topology fronts Postgres with pgbouncer for connection
; pooling. Minimal mode has only 4 Postgres-touching services on a single
; host — direct asyncpg connections to Postgres are fine at this scale,
; and avoiding pgbouncer means one less moving part for a self-host
; security review.

; ── Identity — JWT issuance + validation, on 8002 ───────────────────────
[program:identity]
command=uvicorn services.identity.main:app --host 127.0.0.1 --port 8002
directory=/app
priority=20
autorestart=true
environment=DATABASE_URL="postgresql+asyncpg://identity_user:%(ENV_IDENTITY_DB_PASSWORD)s@postgres:5432/acp_identity"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Registry — agent CRUD + permissions, on 8001 ────────────────────────
[program:registry]
command=uvicorn services.registry.main:app --host 127.0.0.1 --port 8001
directory=/app
priority=21
autorestart=true
environment=DATABASE_URL="postgresql+asyncpg://registry_user:%(ENV_REGISTRY_DB_PASSWORD)s@postgres:5432/acp_registry"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Policy — wraps OPA, exposes /policy/evaluate, on 8003 ───────────────
[program:policy]
command=uvicorn services.policy.main:app --host 127.0.0.1 --port 8003
directory=/app
priority=22
autorestart=true
environment=OPA_URL="http://127.0.0.1:8181"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Behavior — anomaly scoring, on 8007 ─────────────────────────────────
; Redis-only. No DATABASE_URL needed (state lives in `acp:behavior:*`
; keys; nothing is persisted to Postgres).
[program:behavior]
command=uvicorn services.behavior.main:app --host 127.0.0.1 --port 8007
directory=/app
priority=23
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Decision engine — combines policy + behavior signals, on 8010 ──────
[program:decision]
command=uvicorn services.decision.main:app --host 127.0.0.1 --port 8010
directory=/app
priority=24
autorestart=true
environment=POLICY_SERVICE_URL="http://127.0.0.1:8003",BEHAVIOR_SERVICE_URL="http://127.0.0.1:8007"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Audit — chain writer + signer + Merkle roots + bundles, on 8004 ────
[program:audit]
command=uvicorn services.audit.main:app --host 127.0.0.1 --port 8004
directory=/app
priority=25
autorestart=true
environment=DATABASE_URL="postgresql+asyncpg://audit_user:%(ENV_AUDIT_DB_PASSWORD)s@postgres:5432/acp_audit",USAGE_SERVICE_URL="http://127.0.0.1:8006",POLICY_SERVICE_URL="http://127.0.0.1:8003"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Autonomy — escalation contracts + overrides, on 8015 ───────────────
[program:autonomy]
command=uvicorn services.autonomy.main:app --host 127.0.0.1 --port 8015
directory=/app
priority=26
autorestart=true
environment=DATABASE_URL="postgresql+asyncpg://autonomy_user:%(ENV_AUTONOMY_DB_PASSWORD)s@postgres:5432/acp_autonomy"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Gateway — the only externally-reachable process ─────────────────────
; Binds 0.0.0.0:8000. Routes /execute, /demo/groq-agent, /auth/token,
; /agents, /compliance/*, etc. Calls every other inner service via
; localhost. The only thing the outside world ever sees.
[program:gateway]
command=uvicorn services.gateway.main:app --host 0.0.0.0 --port 8000 --workers 2
directory=/app
priority=30
autorestart=true
environment=IDENTITY_SERVICE_URL="http://127.0.0.1:8002",REGISTRY_SERVICE_URL="http://127.0.0.1:8001",POLICY_SERVICE_URL="http://127.0.0.1:8003",DECISION_SERVICE_URL="http://127.0.0.1:8010",AUDIT_SERVICE_URL="http://127.0.0.1:8004",BEHAVIOR_SERVICE_URL="http://127.0.0.1:8007",AUTONOMY_SERVICE_URL="http://127.0.0.1:8015",DATABASE_URL="postgresql+asyncpg://audit_user:%(ENV_AUDIT_DB_PASSWORD)s@postgres:5432/acp_audit",MALLOC_ARENA_MAX="2",PYTHONMALLOC="malloc"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
